# storage-utils/clean-storage-svg-to-react.timer
[Unit]
Description=Delete after 1 days

[Timer]
OnCalendar=daily
Persistent=false

[Install]
WantedBy=timers.target